# config/psiphond-ng.service
[Unit]
Description=Psiphon Next Generation Linux Daemon
Documentation=man:psiphond-ng(8) file:///usr/share/doc/psiphond-ng/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=psiphon
Group=psiphon

# Paths
Environment="CONFIG_FILE=/etc/psiphon/psiphond-ng.conf"
Environment="DATA_DIR=/var/lib/psiphon"
Environment="LOG_FILE=/var/log/psiphon/psiphond-ng.log"

# Executable
ExecStart=/usr/local/bin/psiphond-ng -config $CONFIG_FILE
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=5

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/psiphon /var/log/psiphon
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallFilter=@system-service
SystemCallFilter=~@resources
SystemCallFilter=~@mount
SystemCallFilter=~@reboot
SystemCallFilter=~@privileged
SystemCallErrorNumber=EPERM

# Capabilities needed for TUN mode
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW

# Resource limits
LimitNOFILE=65536
LimitNPROC=512

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
